# mindylu-scheduler.service
[Unit]
Description=MindyLu APScheduler Daemon
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/var/www/mindylu/backend
Environment="PATH=/var/www/mindylu/backend/venv/bin"
ExecStart=/var/www/mindylu/backend/venv/bin/python manage.py run_scheduler
Restart=always

[Install]
WantedBy=multi-user.target
